# scripts/startup/startup-script.service
[Unit]
Description=Startup script

[Service]
Type=simple
ExecStart=/srv/scripts/startup.sh

[Install]
WantedBy=multi-user.target
